# crontab
*/5 * * * * /home/minecraft/pause_minecraft.sh